# com.hirescout.dailyagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this scheduled job -->
    <key>Label</key>
    <string>com.hirescout.dailyagent</string>

    <!-- The command to run: node agent.js -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>HIRESCOUT_PATH/agent.js</string>
    </array>

    <!-- Working directory (same folder as agent.js) -->
    <key>WorkingDirectory</key>
    <string>HIRESCOUT_PATH</string>

    <!-- Run every day at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Log output to agent.log in the hirescout folder -->
    <key>StandardOutPath</key>
    <string>HIRESCOUT_PATH/agent.log</string>
    <key>StandardErrorPath</key>
    <string>HIRESCOUT_PATH/agent.log</string>

    <!-- Start even if the scheduled time was missed (e.g. Mac was asleep) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
